# assets/cert/renew.timer
[Unit]
Description=Renew certificate

[Timer]
Persistent=true
OnCalendar=*-*-* 3:35
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target